# fireplace.service
[Unit]
Description=Fireplace Control Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/fireplace-service
Environment=PATH=/home/pi/.local/bin:/usr/bin:/bin
ExecStart=/home/pi/.local/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
